# debian/limn.upstart
# Limn Server
 
description "Limn"
 
start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]
 
respawn
respawn limit 2 5
 
umask 007
 
kill timeout 300
 
env NODE_ENV=production
env LIMN_PORT=8081
env LIMN_VARDIR=/var/lib/limn
 
setuid limn
setgid limn
 
console output
chdir /usr/lib/limn
exec ./server/server.co >> /var/log/limn/limn.log
